# drakupdaterobot.cron
#!/bin/sh
if [ -f /etc/drakupdatetxt.conf ]; then drakupdaterobot; fi
